# move-original-images.service
[Service]
User = root
Restart = always
Type = forking
ExecStart = /home/scripts/move-original-images.sh
TimeoutSec=0

[Install]
WantedBy=multi-user.target


